# dot_config/systemd/user/gammastep-indicator.service.d/override.conf
[Unit]
ConditionEnvironment=WAYLAND_DISPLAY
StartLimitIntervalSec=60
StartLimitBurst=20
Wants=geoclue2-agent.session

[Service]
Slice=background.slice
RestartSec=3
