# dot_config/systemd/user/pipewire.service.d/50-roc-log.conf
[Service]
# Suppress mod.roc.lib spam (sendto: Network is unreachable when hifipi is offline)
# ROC failures are obvious from missing audio — logs aren't needed for normal operation.
# Remove this file temporarily if debugging ROC issues.
Environment="PIPEWIRE_DEBUG=mod.roc.lib:0"
