# dot_config/systemd/user/hyprsunset.service
[Unit]
Description=Day/night gamma adjustments for Hyprland
Documentation=https://github.com/hyprwm/hyprsunset
OnFailure=notify@%n.service
Requires=graphical-session.target
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=simple
Slice=session.slice
ExecStart=/usr/local/bin/hyprsunset -t 4500
Restart=on-failure
RestartSec=3

[Install]
WantedBy=wlr-session.target
